# docker/Docker.Supervisord.conf
[supervisord]
nodaemon=true

[program:ocpu]
command=apachectl -DFOREGROUND

[program:uwsgi]
command=/usr/local/bin/uwsgi --ini /var/www/Python/uwsgi.ini

[program:nginx]
command=/usr/sbin/nginx